# package/rootfs/datetime_start.service
#
# This file is distributed as part of Xilinx ARM SDK
#
# Copyright (c) 2020 - 2021,  Xilinx, Inc.
# All rights reserved.
#
[Unit]
Description=Date Time start
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=on-failure
User=root
ExecStart=/sbin/datetime_update

[Install]
WantedBy=multi-user.target
